# recipes-setup/gnss-station/files/locupdate.cron
*/2 * * * * root /usr/local/bin/locupdate > /tmp/locupdate.out
