# recipes-setup/gnss-station/files/locupdate.cron
* * * * * root /usr/local/bin/locupdate | tee -a /dev/tty3 > /tmp/locupdate.out
