# recipes-setup/gnss-station/files/rrdcached.cron
5 3 * * * root /etc/init.d/rrdcached restart

