# recipes-setup/gnss-station/files/rrdgraph.cron
*/5 * * * * root /usr/local/bin/rrdgraph > /tmp/rrdgraph.out
