# recipes-setup/gnss-station/files/rtkrcv.cron
8 1 * * * root find /tmp -type f -name "igu*" -mtime +1 -delete
